# amazon-ssm.rc-service
#!/sbin/openrc-run
command="${SSM_BINARY:-/usr/bin/amazon-ssm-agent}"
command_background=true
pidfile="/var/run/amazon-ssm-agent.pid"

rc_ulimit="${SSM_ULIMIT:--c unlimited -n 1048576 -u unlimited}"
SSM_LOGFILE="${SSM_LOGFILE:-/var/log/${RC_SVCNAME}.log}"

depend() {
	need net
}

# Overwriting stop function as the default implementation doesn't wait enough
# for the amazon-ssm-agent to kill child processes before exiting gracefully
stop() {
    ebegin "Stoping amazon-ssm"
    start-stop-daemon --stop \
        --pidfile $pidfile \
        --retry 60
    eend $?
}

